# scripts/com.caspar.yahoo-grab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Fetch live prices from Yahoo Finance every 60 min (IBKR-free fallback) -->
  <key>Label</key>
  <string>com.caspar.yahoo-grab</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/xynkro/Documents/Trading/FinancePWA/.venv/bin/python</string>
    <string>/Users/xynkro/Documents/Trading/FinancePWA/src/yahoo_grab.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/xynkro/Documents/Trading/FinancePWA</string>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>/Users/xynkro/Documents/Trading/FinancePWA/.state/yahoo-grab-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/xynkro/Documents/Trading/FinancePWA/.state/yahoo-grab-stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
